# daemon/claude-relay.service
[Unit]
Description=Claude Code Telegram Relay
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=radxa
WorkingDirectory=/home/radxa/claude-telegram-relay
ExecStart=/home/radxa/.bun/bin/bun run src/relay.ts
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Environment
Environment=HOME=/home/radxa
Environment=PATH=/home/radxa/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=append:/home/radxa/claude-telegram-relay/logs/claude-telegram-relay.log
StandardError=append:/home/radxa/claude-telegram-relay/logs/claude-telegram-relay.error.log

# Resource limits
MemoryMax=500M
CPUQuota=30%

[Install]
WantedBy=multi-user.target
